# com.miapicchietti.syntheticmonitor.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.miapicchietti.syntheticmonitor</string>
    <key>ProgramArguments</key>
    <array>
        <string>/path/to/synthetic-monitor/venv/bin/python3</string>
        <string>/path/to/synthetic-monitor/monitor.py</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/path/to/synthetic-monitor</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>SLACK_WEBHOOK_URL</key>
        <string>YOUR_SLACK_WEBHOOK_URL_HERE</string>
        <key>NETLIFY_AUTH_TOKEN</key>
        <string>YOUR_NETLIFY_PERSONAL_ACCESS_TOKEN_HERE</string>
        <key>NETLIFY_SITE_ID</key>
        <string>YOUR_NETLIFY_SITE_ID_HERE</string>
        <key>NETLIFY_SITE_URL</key>
        <string>YOUR_NETLIFY_SITE_URL_HERE</string>
    </dict>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/path/to/synthetic-monitor/logs/monitor.log</string>
    <key>StandardErrorPath</key>
    <string>/path/to/synthetic-monitor/logs/monitor.err.log</string>
</dict>
</plist>
